# root/etc/services.d/nginx/run
#!/usr/bin/with-contenv sh

#exec ep /etc/nginx/conf.d/default.conf;
exec nginx;
